# test_driven_development_with_python/superlists/deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid doc
chdir /var/www/SITENAME/source

exec ../envs/SITENAME/bin/gunicorn --bind unix:/tmp/SITENAME.socket SITENAME.wsgi:application
